# scripts/com.stock_monitor.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stock_monitor.daemon</string>

    <!-- Calls start_daemon.sh (mirrors StockMonitor-Start Windows task) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tobala/projects/stock_monitor/scripts/start_daemon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tobala/projects/stock_monitor</string>

    <!-- CR-SEC-06: EnvironmentVariables block removed.
         LINE_CHANNEL_ACCESS_TOKEN is retrieved from macOS Keychain
         by start_daemon.sh at runtime via:
           security find-generic-password -s stock_monitor -a LINE_TOKEN -w
         One-time setup:
           security add-generic-password -s stock_monitor -a LINE_TOKEN    -w "<token>"
           security add-generic-password -s stock_monitor -a LINE_GROUP_ID -w "<group_id>"
    -->

    <key>RunAtLoad</key>
    <false/>

    <!-- Weekdays 08:50 — same timing as Windows StockMonitor-Start task -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_start.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_start.log</string>
</dict>
</plist>
