# scripts/com.stock_monitor.stop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stock_monitor.stop</string>

    <!-- Calls stop_daemon.sh (mirrors StockMonitor-Stop Windows task) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tobala/projects/stock_monitor/scripts/stop_daemon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tobala/projects/stock_monitor</string>

    <key>RunAtLoad</key>
    <false/>

    <!-- Weekdays 14:30 — same timing as Windows StockMonitor-Stop task -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_stop.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_stop.log</string>
</dict>
</plist>
