# scripts/com.stock_monitor.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stock_monitor.daemon</string>

    <!-- Calls start_daemon.sh (mirrors StockMonitor-Start Windows task) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tobala/projects/stock_monitor/scripts/start_daemon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tobala/projects/stock_monitor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LINE_CHANNEL_ACCESS_TOKEN</key>
        <string>REPLACE_WITH_YOUR_TOKEN</string>
        <key>LINE_TO_GROUP_ID</key>
        <string>REPLACE_WITH_YOUR_GROUP_ID</string>
        <key>APP_TIMEZONE</key>
        <string>Asia/Taipei</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- Weekdays 08:50 — same timing as Windows StockMonitor-Start task -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>50</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_start.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tobala/projects/stock_monitor/logs/launchd_start.log</string>
</dict>
</plist>
